# graphite-server/packet/etc/rc.local
#!/bin/bash -x

function custom() {
	. /etc/profile.d/ruby.sh
	cd /root/chef/cookbooks && git pull && chef-solo -c /home/ubuntu/solo.rb -j /home/ubuntu/node.json
}

custom >> /var/log/rc.local.log 2>&1