# graphite-server/packet/etc/rc.local
cd /home/ubuntu/sandbox_cookbooks
git update
chef-solo -c /home/ubuntu/solo.rb -j /home/ubuntu/node.json